# liveearth.timer
[Unit]
Description=Update Live Earth Background

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h 

[Install]
WantedBy=timers.target
