# liveearth.service
[Unit]
Description=Live Earth Background

[Service]
ExecStart=%h/.local/bin/bg.sh

[Install]
WantedBy=default.target